# services.d/lb-config-gen/run
#!/usr/bin/with-contenv bash

exec /lb-config-gen.py

